# bamf/bamfwatcher.service
[Unit]
Description=Update Bamf index for Unity shell

[Service]
Type=simple
ExecStart=/usr/bin/bamfwatcher
Restart=always

[Install]
WantedBy=multi-user.target
